# confs/deploy_results.conf
start on runlevel [2345]
stop on runlevel [!2345]

respawn
setuid ubuntu
setgid ubuntu

script
    . /etc/environment
    /bin/bash /home/ubuntu/apps/elections14/repository/run_on_server.sh fab $DEPLOYMENT_TARGET daemons.deploy_results &> {{ SERVER_LOG_PATH }}/deploy_results.log
end script

