# confs/rotate_slide.conf
start on runlevel [2345]
stop on runlevel [!2345]

respawn

script
    . /etc/environment
    /bin/bash /home/ubuntu/apps/elections14/repository/run_on_server.sh fab $DEPLOYMENT_TARGET daemons.rotate_slide >/var/log/elections14_rotate_slide.log
end script
